# docker/crontab
* * * * * php /var/www/yii grabber >> /var/log/cron.log 2>&1
